# deploy/password-pdf-generator.service.example
[Unit]
Description=Password PDF Generator API
After=network.target

[Service]
User=passwordpdf
Group=passwordpdf
WorkingDirectory=/opt/services/password-pdf-generator/app
EnvironmentFile=/opt/services/password-pdf-generator/config/password-pdf-generator.env
Environment=WIFI_PDF_CONFIG_PATH=/opt/services/password-pdf-generator/config/brand_settings.json
Environment=WIFI_PDF_LOG_DIR=/opt/services/password-pdf-generator/logs
Environment=PATH=/opt/services/password-pdf-generator/.venv/bin
ExecStart=/opt/services/password-pdf-generator/.venv/bin/uvicorn wifi_pdf.api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
